# deploy/jin-tian-ni-shou-le-ma.service
[Unit]
Description=Jin Tian Ni Shou Le Ma weight tracking site
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/jin-tian-ni-shou-le-ma
ExecStart=/usr/bin/node /opt/jin-tian-ni-shou-le-ma/server.mjs
Environment=NODE_ENV=production
EnvironmentFile=-/etc/jin-tian-ni-shou-le-ma/admin.env
Environment=DATA_DIR=/opt/jin-tian-ni-shou-le-ma/data
Environment=ALLOWED_HOSTS=furby.top,www.furby.top,127.0.0.1
Environment=SSL_KEY_PATH=/etc/jin-tian-ni-shou-le-ma/furby.top.key
Environment=SSL_CERT_PATH=/etc/jin-tian-ni-shou-le-ma/furby.top_bundle.pem
Environment=PORT=443
Environment=HTTP_PORT=80
Environment=TEST_PROXY_PATH=/test
Environment=TEST_PROXY_TARGET=http://127.0.0.1:3101
Restart=always
RestartSec=3
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
ReadWritePaths=/opt/jin-tian-ni-shou-le-ma/data
LimitNOFILE=8192
TasksMax=128
MemoryMax=768M

[Install]
WantedBy=multi-user.target
